# systemd/birdwatch.socket
[Unit]
Description=HTTP Socket for Birdwatch

[Socket]
ListenStream=[::1]:8077
Accept=yes

[Install]
WantedBy=sockets.target